# crontab
0 5,12,9 * * * python /getSchedule.py >> /logfile
